# config/supervisor.ini
[supervisord]
nodaemon=true
user = root

[program:echo]
command=trace-test
user = root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=LANG=en_US.UTF-8,LC_ALL=en_US.UTF-8

[program:envoy]
command=func-e run -c /opt/trace-test/etc/envoy.conf.yaml --service-cluster c0 --service-node n0 --log-level %(ENV_ENVOY_LOG_LEVEL)s
user = root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
stopsignal=QUIT